# luci-app-lucky/root/etc/init.d/lucky
#!/bin/sh /etc/rc.common
START=99
STOP=15
USE_PROCD=1

CRON_FILE=/etc/crontabs/root
CRON_TAG=luci-app-lucky-autoupdate
UPDATE_SCRIPT=/usr/libexec/lucky-update.sh
LOGFILE=/tmp/lucky_update/lucky-start.log

log_msg() { logger -t lucky -p daemon.info "$1"; }

init_config() {
    local v key var def kv rest

    for kv in \
        "enabled=ENABLED=0" \
        "binpath=PROG=/usr/bin/lucky" \
        "configdir=CONFDIR=/etc/config/lucky.daji" \
        "port=PORT=16601" \
        "safe=SAFE=" \
        "delay=DELAY=60" \
        "respawn_threshold=RESPAWN_THRESHOLD=3600" \
        "respawn_timeout=RESPAWN_TIMEOUT=30" \
        "respawn_retry=RESPAWN_RETRY=5" \
        "auto_update=AUTO_UPDATE=0" \
        "update_interval=UPDATE_INTERVAL=7"
    do
        key="${kv%%=*}"
        rest="${kv#*=}"
        var="${rest%%=*}"
        def="${rest#*=}"
        v="$(uci -q get lucky.lucky.$key)"
        eval "$var=\"\${v:-$def}\""
    done

    mkdir -p "$CONFDIR" "${LOGFILE%/*}"
}

sync_cron() {
    init_config

    local current
    current=$(grep "$CRON_TAG" "$CRON_FILE" 2>/dev/null)

    if [ "$AUTO_UPDATE" = "1" ]; then
        case "$current" in
            *"*/${UPDATE_INTERVAL} * * $UPDATE_SCRIPT"*) return 0 ;;
        esac
        sed -i "/${CRON_TAG}/d" "$CRON_FILE" 2>/dev/null
        echo "0 3 */${UPDATE_INTERVAL} * * $UPDATE_SCRIPT #${CRON_TAG}" >> "$CRON_FILE"
        log_msg "Auto update cron set: every ${UPDATE_INTERVAL} day(s)"
    else
        [ -z "$current" ] && return 0
        sed -i "/${CRON_TAG}/d" "$CRON_FILE" 2>/dev/null
        log_msg "Auto update cron removed."
    fi

    /etc/init.d/cron restart 2>/dev/null
}

start_service() {
    init_config

    if [ ! -x "$PROG" ]; then
        uci -q set lucky.lucky.installed_version='Unknown' 2>/dev/null
        uci -q commit lucky 2>/dev/null
        [ -x "$UPDATE_SCRIPT" ] && {
            log_msg "Lucky binary not found, triggering download..."
            setsid sh -c "$UPDATE_SCRIPT" >/dev/null 2>&1 &
        }
        return 0
    fi

    [ "$ENABLED" = "1" ] || { log_msg "Lucky is disabled, skipping start."; return 0; }

    [ "$(awk -F. '{print $1}' /proc/uptime)" -lt 120 ] && {
        log_msg "Delaying start by ${DELAY}s..."
        sleep "$DELAY"
    }

    "$PROG" -setconf -key AdminWebListenPort -value "$PORT" -cd "$CONFDIR" >/dev/null 2>&1

    if [ -z "$SAFE" ]; then
        "$PROG" -rCancelSafeURL -cd "$CONFDIR" >/dev/null 2>&1
    else
        "$PROG" -setconf -key SafeURL -value "$SAFE" -cd "$CONFDIR" >/dev/null 2>&1
    fi

    : > "$LOGFILE"

    procd_open_instance
    procd_set_param command /bin/sh -c \
        "exec \"$PROG\" -cd \"$CONFDIR\" >> \"$LOGFILE\" 2>&1"
    procd_set_param respawn "$RESPAWN_THRESHOLD" "$RESPAWN_TIMEOUT" "$RESPAWN_RETRY"
    procd_close_instance

    log_msg "Lucky started (port=${PORT} configdir=${CONFDIR})"
}

stop_service() { log_msg "Lucky stopped."; }

reload_service() { sync_cron; }

service_triggers() { procd_add_reload_trigger lucky; }